# ops/crontab
# m h  dom mon dow   command
30 18 * * * /bin/sh /var/maltaa/backup_main.sh
40 18 * * * /bin/sh /var/maltaa/backup_sync.sh